# supervisor/conf.d/dev-reportcard.conf
[program:dev-reportcard]
command=/srv/dev-reportcard.wmflabs.org/limn/src/server/server.co
directory=/srv/dev-reportcard.wmflabs.org/limn
user=www
environment=NODE_ENV=dev,LIMN_PORT=8081

