# supervisor/conf.d/gerrit-stats.conf
[program:gerrit-stats]
command=/srv/gerrit-stats.wmflabs.org/limn/src/server/server.co
directory=/srv/gerrit-stats.wmflabs.org/limn
user=www
environment=NODE_ENV=development,LIMN_PORT=8082

