# supervisor/conf.d/global-dev.conf
[program:global-dev]
command=/srv/global-dev.wmflabs.org/limn/src/server/server.co
directory=/srv/global-dev.wmflabs.org/limn
user=www
environment=NODE_ENV=development,LIMN_PORT=8085

